# infra/systemd/nhms-compute-compose.service
# Example systemd unit for node 22 compute/control Docker Compose.
#
# Operator edits required before install:
# - Replace /opt/SHUD-NWM with the absolute repository checkout path on node 22.
# - Create infra/env/compute.env from infra/env/compute.example with real
#   compute-api DB credentials, DB-free scheduler backend/path selectors,
#   workspace/object-store/published/scheduler roots, Basins/model paths,
#   source/model filters, and image tag. The scheduler-once compose service must
#   not receive DATABASE_URL; keep that value only for compute-api / rollback.
# - Ensure the docker compose plugin is installed at /usr/bin/docker or update
#   ExecStart/ExecStop/ExecReload to the local docker binary path.
# - Before enable/start/reload/restart, run
#   scripts/validate_two_node_docker_source_trust.py --role compute from the
#   checkout root and record its docker-security evidence. It must pass before
#   systemd consumes WorkingDirectory, compose, env, or unit sources. Docker
#   access is root-equivalent.
# - Business scheduler validation is the manual scheduler-once compose lane:
#   docker compose --env-file env/compute.env -f compose.compute.yml run --rm scheduler-once
#   It runs nhms-pipeline plan-production --plan without --workspace-root,
#   --lock-path, or --evidence-dir. Those roots must come from compute.env;
#   explicit root flags are diagnostic compatibility only. Its runtime env must
#   show NHMS_SCHEDULER_DB_FREE_REQUIRED=true, every scheduler backend selector
#   set to file, and DATABASE_URL absent.
# - Optional bounded scheduler timer lane uses the same compute.env roots and
#   filters. Example companion units:
#     nhms-compute-scheduler.service:
#       Type=oneshot
#       WorkingDirectory=/opt/SHUD-NWM/infra
#       ExecStart=/usr/bin/docker compose --env-file env/compute.env -f compose.compute.yml run --rm scheduler-once uv run nhms-pipeline plan-production --continuous --interval-seconds ${NHMS_SCHEDULER_INTERVAL_SECONDS} --max-passes ${NHMS_SCHEDULER_MAX_PASSES} --max-cycles-per-source ${NHMS_SCHEDULER_MAX_CYCLES_PER_SOURCE} --source gfs --model-id basins_qhh_shud --basin-id basins_qhh --plan
#     nhms-compute-scheduler.timer:
#       OnBootSec=2min
#       OnActiveSec=30s
#       OnUnitActiveSec=5min
#       Persistent=false
#       Unit=nhms-compute-scheduler.service
#   Disable with:
#     systemctl disable --now nhms-compute-scheduler.timer
#     systemctl reset-failed nhms-compute-scheduler.service
#   This timer lane is the production continuous daemon
#   (nhms-pipeline plan-production --continuous -> services/orchestrator/scheduler.py),
#   submitting via the Slurm gateway (SLURM_GATEWAY_* in compute.env); the legacy
#   run_qhh_* scripts are diagnostic-only fallbacks (M24 §5 / #293). Live daemon
#   receipt: artifacts/m24/m24-daemon-5880d09/ (lease_nfs.json / grib_preflight.json /
#   daemon_pass.json).

[Unit]
Description=NHMS node 22 compute/control compose
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/SHUD-NWM/infra
ExecStart=/usr/bin/docker compose --env-file env/compute.env -f compose.compute.yml up -d
ExecStop=/usr/bin/docker compose --env-file env/compute.env -f compose.compute.yml down
ExecReload=/usr/bin/docker compose --env-file env/compute.env -f compose.compute.yml up -d
RemainAfterExit=yes
TimeoutStartSec=0
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
