# infra/systemd/nhms-compute-compose.service
# Example systemd unit for node 22 compute/control Docker Compose.
#
# Operator edits required before install:
# - Replace /opt/SHUD-NWM with the absolute repository checkout path on node 22.
# - Create infra/env/compute.env from infra/env/compute.example with real
#   writer DB credentials, workspace paths, Basins/model paths, and image tag.
# - Ensure the docker compose plugin is installed at /usr/bin/docker or update
#   ExecStart/ExecStop/ExecReload to the local docker binary path.
# - Before enable/start/reload/restart, record the checkout trust preflight from
#   infra/README.two-node-docker.md: every WorkingDirectory path component and
#   unit/compose/env source must be owned by root or a trusted Docker operator
#   and must not be group/world-writable. Docker access is root-equivalent.

[Unit]
Description=NHMS node 22 compute/control compose
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/SHUD-NWM/infra
ExecStart=/usr/bin/docker compose --env-file env/compute.env -f compose.compute.yml up -d
ExecStop=/usr/bin/docker compose --env-file env/compute.env -f compose.compute.yml down
ExecReload=/usr/bin/docker compose --env-file env/compute.env -f compose.compute.yml up -d
RemainAfterExit=yes
TimeoutStartSec=0
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
