# infra/systemd/nhms-node27-download.timer
[Unit]
Description=Run NHMS node-27 raw GFS/IFS source download every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=false
Unit=nhms-node27-download.service

[Install]
WantedBy=timers.target
