# infra/systemd/nhms-scheduler-file-provider-refresh.service
[Unit]
Description=NHMS node-22 DB-free scheduler file-provider refresh
Before=nhms-compute-scheduler.service

[Service]
Type=oneshot
ExecCondition=/bin/sh -c '! /usr/bin/systemctl --user is-active --quiet nhms-compute-scheduler.service'
WorkingDirectory=/scratch/frd_muziyao/NWM
EnvironmentFile=/scratch/frd_muziyao/NWM/infra/env/compute.scheduler-provider-refresh.env
Environment=PATH=/scratch/frd_muziyao/NWM/.venv/bin:/usr/local/bin:/usr/bin:/bin
UnsetEnvironment=DATABASE_URL PIPELINE_DATABASE_URL PGAPPNAME PGCHANNELBINDING PGCLIENTENCODING PGCONNECT_TIMEOUT PGDATABASE PGDATESTYLE PGGEQO PGGSSDELEGATION PGGSSENCMODE PGGSSLIB PGHOST PGHOSTADDR PGKRBSRVNAME PGLOADBALANCEHOSTS PGLOCALEDIR PGMAXPROTOCOLVERSION PGMINPROTOCOLVERSION PGOPTIONS PGPASSFILE PGPASSWORD PGPORT PGREQUIREAUTH PGREQUIREPEER PGREQUIRESSL PGSERVICE PGSERVICEFILE PGSSLCERT PGSSLCERTMODE PGSSLCOMPRESSION PGSSLCRL PGSSLCRLDIR PGSSLKEY PGSSLMAXPROTOCOLVERSION PGSSLMINPROTOCOLVERSION PGSSLMODE PGSSLNEGOTIATION PGSSLROOTCERT PGSSLSNI PGSSL_CERT_FILE PGSSL_KEY_FILE PGSSL_ROOT_CERT_FILE PGSYSCONFDIR PGTARGETSESSIONATTRS PGTZ PGUSER
ExecStart=/scratch/frd_muziyao/NWM/scripts/scheduler_file_provider_refresh_once.sh
StandardOutput=journal
StandardError=journal
# node-22 user-systemd makes the filesystem root inaccessible inside a
# PrivateTmp mount namespace. The provider's no-follow verifier must open the
# root and walk each absolute path component, so this unit cannot use PrivateTmp.
UMask=0077
TimeoutStartSec=7200
