# infra/systemd/nhms-scheduler-evidence-retention.service
[Unit]
Description=NHMS node-22 scheduler evidence retention

[Service]
Type=oneshot
WorkingDirectory=/scratch/frd_muziyao/NWM
EnvironmentFile=/scratch/frd_muziyao/NWM/infra/env/compute.scheduler-dbfree.env
ExecStart=/usr/bin/env uv run python scripts/node22_scheduler_evidence_retention.py
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600
