# infra/systemd/nhms-node27-timeseries-retention.timer
[Unit]
Description=Run NHMS node-27 timeseries retention daily

[Timer]
OnCalendar=*-*-* 05:15:00 UTC
Persistent=true
Unit=nhms-node27-timeseries-retention.service

[Install]
WantedBy=timers.target
