# infra/systemd/nhms-scheduler-file-provider-refresh.timer
[Unit]
Description=Refresh NHMS scheduler file providers daily

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
RandomizedDelaySec=30m
Persistent=false
Unit=nhms-scheduler-file-provider-refresh.service

[Install]
WantedBy=default.target
